# src/config/os/etc/systemd/system/omega-frontend.service
[Unit]
Description=OMEGA CORE Frontend (Vite)
After=network.target

[Service]
Type=simple
User=mth
WorkingDirectory=/OMEGA_CORE/frontend
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
