# infra/vps/systemd/omega-core-anti-heroin.timer
[Unit]
Description=Täglicher OMEGA CORE Anti-Heroin-Scan

[Timer]
OnCalendar=daily
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
